# autostart/.config/autostart/Touchegg.desktop
[Desktop Entry]
Encoding=UTF-8
Name=touchegg
Comment=TouchEgg Gestures
Exec=touchegg &
Type=Application
